# compose-recipes/podman/s6-rc.d/deploy-client-extensions/run
#!/bin/bash
echo "Starting Kopf Operator..."
exec kopf run /src/deploy-client-extensions.py --verbose --standalone