# compose-recipes/podman/s6-rc.d/apply-crd/up
#!/bin/sh
exec kubectl apply -f /src/liferayextensions.yaml
