# compose-recipes/podman/s6-rc.d/watch-configmaps/run
#!/bin/bash
echo "Starting Kopf Operator..."
exec kopf run /src/watch-configmaps.py --verbose --standalone