# compose-recipes/k8s/s6-rc.d/k8s-proxy/run
#!/bin/sh
exec kubectl proxy --kubeconfig=/var/lib/rancher/k3s/server/cred/admin.kubeconfig  --port=8001 --address=0.0.0.0 --accept-hosts='.*'
